# web/root/etc/services.d/web/run
#!/bin/sh

# Be sure to set our cwd correctly..
cd /app
# start our node.js application
exec node server/server.js
